# systemd/raspi-monitor-prune.timer
[Unit]
Description=Daily raspi_monitor metrics prune (30-day retention)
[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
[Install]
WantedBy=timers.target
